# udev/99-monitor-hotplug.rules
# Monitor hotplug rule for hypr-mon-switch
# Triggers monitor layout changes when displays are connected/disconnected
ACTION=="change", SUBSYSTEM=="drm", RUN+="/etc/acpi/monitor-hotplug-config.sh"
